# DEBIAN.template/postinst
#!/bin/bash
set -ex

ls -lah /usr/bin/process-exporter

systemctl daemon-reload
systemctl enable td-process-exporter
systemctl start td-process-exporter
