# DEBIAN.template/postinst
#!/bin/bash
set -ex

#systemctl daemon-reload
#systemctl enable process-exporter
#systemctl start process-exporter
service process-exporter start
