# DEBIAN.template/postinst
#!/bin/bash
set -ex

ls -lah /usr/bin/process-exporter

systemctl daemon-reload
systemctl enable process-exporter
systemctl start process-exporter
