# etc/systemd/system/process-exporter.service
[Unit]
Description = Process exporter
After = network.target

[Service]
ExecStart = /usr/bin/process-exporter -config.path /etc/process-exporter/conf.yaml -web.listen-address 127.0.0.1:7756 -web.telemetry-path /psmon/metrics

[Install]
WantedBy = multi-user.target
